# ansible/roles/bichon/templates/bichon-archive.timer.j2
[Unit]
Description=Run bichon email archive hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10min
Persistent=true
Unit=bichon-archive.service

[Install]
WantedBy=timers.target
