# ansible/roles/hermes/templates/hermes-gateway.service.j2
[Unit]
Description=Hermes Agent Gateway - Personal AI Assistant
After=network.target

[Service]
Type=simple
WorkingDirectory={{ hermes_config_dir }}
Environment="PATH={{ hermes_home }}/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile={{ hermes_config_dir }}/.env
ExecStart={{ hermes_home }}/.local/bin/hermes gateway
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
