# ansible/roles/baikal/templates/baikal-busy-sync.timer.j2
[Unit]
Description=Baikal busy feed sync timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
