# ansible/roles/vdirsyncer/templates/vdirsyncer.timer.j2
[Unit]
Description=vdirsyncer iCloud-to-Baikal calendar sync timer

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ vdirsyncer_sync_interval }}
AccuracySec=30s

[Install]
WantedBy=timers.target
